# continuous_integration/docker/hadoop/files/etc/supervisord.d/fixuid.conf
[program:fixuid]
command=/root/fixuid.sh
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/fixuid.log
autostart=true
autorestart=false
